# Backend/sandbox/supervisord.conf
; Minimal "execution sandbox" runtime: Xvfb -> x11vnc (view-only) -> websockify/noVNC
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac +extension RANDR
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xfce]
command=/bin/bash -lc "dbus-launch --exit-with-session startxfce4"
environment=DISPLAY=%(ENV_DISPLAY)s,HOME=/home/pwuser,USER=pwuser
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:chrome]
command=/bin/bash -lc "dbus-launch /home/pwuser/sandbox/launch_chrome.sh"
environment=DISPLAY=%(ENV_DISPLAY)s,HOME=/home/pwuser,USER=pwuser,CHROME_URL=%(ENV_CHROME_URL)s
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20

[program:x11vnc]
command=/usr/bin/x11vnc -display %(ENV_DISPLAY)s -rfbport 5900 -viewonly -shared -forever -nopw -listen 0.0.0.0 -noxdamage
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; Ubuntu package installs noVNC under /usr/share/novnc/
; websockify should serve static assets and bridge ws://:6080 -> tcp://127.0.0.1:5900
[program:novnc]
command=/usr/bin/websockify --web=/home/pwuser/sandbox/novnc-web/ 0.0.0.0:6080 127.0.0.1:5900
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

